# files/MedicalCoderSwarm-main/api/systemd/mcs-docker.service
# derived from https://github.com/encode/uvicorn/issues/678
# dervied from https://blog.container-solutions.com/running-docker-containers-with-systemd
[Unit]
Description=mcs
After=docker.service
#Required=docker.service

[Service]
EnvironmentFile=/var/run/mcs/secrets/env
RestartSec=10
TimeoutStartSec=0
Restart=always
ExecStartPre=-/usr/bin/docker stop %n
ExecStartPre=-/usr/bin/docker rm %n
ExecStartPre=/usr/bin/docker pull h4ckermike/mcs-api:experimental
#ExecStart=/usr/bin/docker run --rm --name %n h4ckermike/mcs-api:experimental
ExecStart=/usr/bin/docker run -p 8000:8000 -w /var/mcs/agent_workspace/ --mount type=bind,source=/opt/mcs,target=/opt/mcs  -e WORKSPACE_DIR=/var/mcs/agent_workspace/ --rm --name %n  h4ckermike/mcs-api:experimental /usr/bin/unbuffer /var/mcs/agent_workspace/.venv/bin/uvicorn --proxy-headers --forwarded-allow-ips='*'--workers=4 --host 0.0.0.0 --port=8000 --reload-delay=30 --app-dir /opt/mcs/api main:create_app
# --network host
# leave these logs for now because they are
# FIXME: update cloudwatch logs
StandardOutput=file:/var/log/swarms_systemd.log
StandardError=file:/var/log/swarms_systemd.log
ExecReload=/bin/kill -HUP ${MAINPID}

Restart=always

[Install]
WantedBy=multi-user.target

